# 2025/2025-05/defendnot/sigma/win_security_defendnot_default_scheduled_task_creation.yml

title: Defendnot Scheduled Task Creation - Security
id: 3156aabd-474b-408e-8fc5-f9a92e026cd1
status: experimental
description: Detects creation of scheduled tasks potentially used for autorun persistence (Event ID 4698), such as those created by defendnot or similar tooling.
references:
    - https://github.com/es3n1n/defendnot
    - https://www.huntress.com/blog/defendnot-detecting-malicious-security-product-bypass-techniques
author: Matt Anderson
date: 2025-05-22
logsource:
    product: windows
    service: security
detection:
    selection:
        EventID: 4698
        TaskName|contains:
            - "defendnot"
            - "autorun"
            - "Loader"
        CommandLine|contains:
            - "defendnot"
            - "--from-autorun"
            - "defendnot-loader"
            - "\\defendnot.dll"
    condition: selection
fields:
    - SubjectUserName
    - TaskName
    - CommandLine
    - ProcessId
    - ComputerName
falsepositives:
    - Legitimate admin tasks with similar names or command lines
level: high
tags:
    - attack.persistence
    - attack.t1053.005